# package/base-files/files/etc/uci-defaults/disable-openwrt-firewall
#!/bin/sh
/etc/init.d/firewall stop
/etc/init.d/firewall disable
exit 0
