# package/base-files/files/etc/uci-defaults/disable-polipo
#!/bin/sh
/etc/init.d/polipo stop
/etc/init.d/polipo disable
exit 0
