# package/base-files/files/etc/uci-defaults/disable-openwrt-qos
#!/bin/sh
/etc/init.d/qos stop
/etc/init.d/qos disable
exit 0
