# package/base-files/files/etc/uci-defaults/disable-dhcp
#!/bin/sh
/etc/init.d/dnsmasq stop
/etc/init.d/dnsmasq disable
exit 0
