# package/base-files/files/etc/uci-defaults/nnx-precfg
#!/bin/sh
/etc/init.d/ninux start
exit 0
