# package/base-files/files/etc/uci-defaults/olsrd-policy-routing-ninux
uci set olsrd.@olsrd[0].RtTable=111
uci set olsrd.@olsrd[0].RtTableDefault=112
uci commit
